# images/rc.local
#!/bin/sh -e
service mysql start
service apache2 start
exit 0
